# roles/monitoring/files/etc/cron.d/check_cron
# This file is autogenerate by ansible
# Do not edit, changes will be overwritten
*/5 * * * *       root    touch /tmp/check_cron
