# ansible/roles/proxy_server/files/etc/cron.daily/clamav-freschlam
# This file is autogenerate by ansible
# Do not edit, changes will be overwritten

#!/bin/sh

freshclam --quiet --no-warnings
